# etc/service/statsd/run
#!/bin/sh
exec /usr/bin/setuidgid www-data /usr/local/bin/node /opt/statsd/stats.js /opt/statsd/localConfig.js
