# etc/service/carbon-cache/run
#!/bin/sh
exec /usr/bin/setuidgid www-data /opt/graphite/bin/carbon-cache.py --debug start
